# data/files/img_reader.timer
[Unit]
Description= runs every hour

[Timer]
OnCalendar=hourly
Unit=img_reader.service
Persistent=true

[Install]
WantedBy=default.target
