# data/files/send_mail.timer
[Unit]
Description= runs once a day a timer

[Timer]
OnCalendar=*-*-* 09:00:00
Unit=send_mail.service
Persistent=true

[Install]
WantedBy=default.target
